# etc/supervisor_uwsgi.ini
[program:uwsgi]
directory = /root
command = /usr/local/bin/uwsgi --emperor "/etc/uwsgi/apps-enabled/*" --die-on-term --master --catch-exceptions
autostart   = true
autorestart = true